# var-mnt-bigdrive.mount
[Unit]
Description=Secondary drive mount

[Mount]
What=/dev/sda1
Where=/var/mnt/bigdrive

[Install]
WantedBy=multi-user.target
